# schedule/elixir_metric.timer
[Unit]
Description=elixir metrics

[Timer]
OnCalendar=*-*-* *:0/2:00
Persistent=true

[Install]
WantedBy=timers.target
